# raspberrypi/tools/mirai.service
[Unit]
Description=Mirai
After=network.target

[Service]
WorkingDirectory=/home/pi/mirai
ExecStart=/home/pi/mirai/app-linux
Restart=on-failure
User=pi

[Install]
WantedBy=multi-user.target
